# site-cookbooks/consul/templates/default/consul-watch.service.erb
[Unit]
Description=consul-watch agent
Requires=consul.service

[Service]
EnvironmentFile=-/etc/sysconfig/consul-watch
Environment=GOMAXPROCS=2
Restart=on-failure
ExecStart=/usr/local/bin/consul watch -type event -name <%= @role %> stretcher
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
