# site-cookbooks/unicorn/templates/default/unicorn.service.erb
[Unit]
Description=unicorn

[Service]
User=app
Group=app
WorkingDirectory=/var/www/current
PIDFile=/var/www/current/tmp/pids/unicorn.pids
ExecStart=/usr/local/bin/bundle exec unicorn -c config/unicorn/<%= @rails_env %>.rb -E <%= @rails_env %>
ExecReload=/bin/kill -USR2 $MAINPID
ExecStop=/bin/kill -QUIT $MAINPID
KillSignal=SIGINT

[Install]
